# rails/package/application-workers.target
[Unit]
Description=[**Application**] Workers
Requires=[**application**].target
After=network.target

[Install]
WantedBy=[**application**].target
