# rails/package/application.target
[Unit]
Description=[**Application**] Application

[Install]
WantedBy=multi-user.target
